# misc/frame-broker.service
[Unit]
Description=Frame Broker Service

[Service]
ExecStart=/usr/bin/node /opt/frame-broker/dist/index.js

[Install]
WantedBy=multi-user.target